# test/docker/logger/systemd-journal-gatewayd.socket
[Unit]
Description=Journal Gateway Service Socket

[Socket]
ListenStream=/run/systemd-journal-gatewayd.sock

[Install]
WantedBy=sockets.target
